# packaging/simpleshortener_0.1-b2_amd64/DEBIAN/postinst
psql -U postgres -c "CREATE DATABASE simpleshortener"
psql -U postgres -c "CREATE USER simpleshortener PASSWORD 'simpleshortener'"
systemctl enable simpleshortener.service
systemctl start simpleshortener.service
echo "Configuration files placed in /etc/simpleshortener/, please edit them as needed"