# packaging/simpleshortener_0.1-b2_amd64/DEBIAN/postinst
systemctl enable simpleshortener.service
systemctl start simpleshortener.service
echo "Configuration files placed in /etc/simpleshortener/, please edit them as needed"